# .config/systemd/user/battery.service
[Unit]
Description="Warn user of low battery"

[Service]
Type=oneshot
ExecStart=%h/.local/share/scripts/battery-warning.sh
