# .config/systemd/user/mail.service
[Unit]
Description="Fetch new mail."

[Service]
Type=oneshot
ExecStart=/usr/bin/env offlineimap
